# deploy/ubuntu/omsg-api.service
[Unit]
Description=oMsg FastAPI
After=network.target

[Service]
Type=simple
User=omsg
Group=omsg
WorkingDirectory=/opt/omsg
EnvironmentFile=/etc/omsg/api.env
ExecStart=/opt/omsg/venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
